# autostart/google-chrome.desktop
[Desktop Entry]
Name=Chrome
Exec=google-chrome --profile-directory=Default
Terminal=false
Type=Application
StartupNotify=false
