# autostart/nvim.desktop
[Desktop Entry]
Name=nvim
Exec=alacritty --class nvim -e nvim --title nvim
Terminal=false
Type=Application
StartupNotify=false
