# autostart/telegram.desktop
[Desktop Entry]
Name=Talegram
Exec=Telegram
Terminal=false
Type=Application
StartupNotify=false
