# autostart/vifm.desktop
[Desktop Entry]
Name=nvim
Exec=alacritty --class vifm -e vifm --title vifm
Terminal=false
Type=Application
StartupNotify=false
